# recipes-tools/netdev-led/files/netdev-led.service
[Unit]
Description=netdev-led is setting the mode of the network LED.
After=network.target
ConditionPathExists=!/update/end_of_life

[Service]
Type=simple
ExecStart=/usr/bin/netdev-led

[Install]
WantedBy=network.target
