# recipes-core/emcfg/files/em-timesync.timer
[Unit]
Description=Trigger system/RTC time synchronization
DefaultDependencies=no

[Timer]
OnUnitActiveSec=1min

[Install]
WantedBy=sysinit.target
