# recipes-core/emcfg/files/emcfg.service
[Unit]
Description=Transient system configuration
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/emcfg update
