# recipes-tools/empkg/files/em-app-time.target
[Unit]
Description=EM apps which need valid time/date
Requires=em-app-before.target
After=em-app-before.target
