# recipes-tools/empkg/files/em-app.target
[Unit]
Description=EM apps
Requires=basic.target
Conflicts=rescue.service rescue.target
After=basic.target rescue.service rescue.target
