# recipes-tools/empkg/files/em-app-no-time.target
[Unit]
Description=EM apps which don't need valid time/date
Requires=em-app-before.target
After=em-app-before.target
