# recipes-tools/empkg/files/em-app.target
[Unit]
Description=EM apps
Requires=em-app-before.target em-app-time.target em-app-no-time.target
After=em-app-before.target
