# recipes-tools/empkg/files/em-app.target
[Unit]
Description=EM apps
Requires=em-app-before.target
After=em-app-before.target
