# recipes-tools/empkg/files/em-app-before.target
[Unit]
Description=EM apps (before)
Requires=basic.target
Conflicts=rescue.service rescue.target
After=basic.target rescue.service rescue.target
