# recipes-tools/empkg/files/em-eebus-use-cases.target
[Unit]
Description=EM EEBus Target
StopWhenUnneeded=true
